# github-scrapy/webapp/ops/deploy/github.ini
[program:github]
command = /alidata1/admin/.virtualenvs/github/bin/gunicorn app:app --config ops/deploy/gunicorn_cfg.py
directory = /alidata1/admin/apps/github/github/webapp
user = admin

[program:verify_proxy]
command = /alidata1/admin/.virtualenvs/github/bin/python scrapy_model.py
directory = /alidata1/admin/apps/github/github/webapp
user = admin

[program:proxy]
command = /alidata1/admin/.virtualenvs/github/bin/python proxy.py --hostname 0.0.0.0
directory = /alidata1/admin/apps/github/github/proxy
user = admin

[program:Xvfb]
command = Xvfb :99 -ac -screen 0 1280x1024x24
user = admin